# dot_config/systemd/user/elephant.service
[Unit]
Description=Elephant launcher backend (Walker provider service)
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target
ConditionEnvironment=XDG_CURRENT_DESKTOP=hyprland:Hyprland

[Service]
Type=simple
ExecStart=/usr/bin/elephant
Restart=on-failure
RestartSec=2
Slice=session.slice

[Install]
WantedBy=graphical-session.target
