# dot_config/systemd/user/mail-sync.service
[Unit]
Description=Sync Gmail via lieer (notmuch + 1Password-backed OAuth)
Documentation=file://%h/Documents/docs/notmuch-gmailieer-setup.md
After=graphical-session.target

[Service]
Type=oneshot
# 1Password's op CLI talks to the desktop app over DBus for biometric/FIDO
# unlock; the desktop session must be active. The mail-sync wrapper handles
# its own tmpfs cleanup via traps.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=%h/.local/bin/mail-sync
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=5

# Stay polite: if the user's busy, don't fight for CPU.
CPUWeight=50
