# dot_config/systemd/user/mail-sync.service
[Unit]
Description=Sync Gmail via lieer (notmuch + TPM2-sealed OAuth)
Documentation=file://%h/Documents/docs/notmuch-gmailieer-setup.md
After=graphical-session.target

[Service]
Type=oneshot
# mail-sync unseals the OAuth credentials from a TPM2-bound blob into tmpfs
# via systemd-creds. No interactive unlock is required and no desktop session
# dependency exists — the timer can fire entirely unattended. The wrapper
# handles its own tmpfs cleanup via traps.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=%h/.local/bin/mail-sync
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=5

# Stay polite: if the user's busy, don't fight for CPU.
CPUWeight=50
